# installer/macos/com.minnowvpn.vpn-service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique identifier for the service -->
    <key>Label</key>
    <string>com.minnowvpn.vpn-service</string>

    <!-- Program and arguments -->
    <key>ProgramArguments</key>
    <array>
        <string>/Library/PrivilegedHelperTools/minnowvpn-service</string>
        <string>--daemon</string>
    </array>

    <!-- Run at system boot -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart policy -->
    <key>KeepAlive</key>
    <dict>
        <!-- Restart if the process crashes (non-zero exit) -->
        <key>SuccessfulExit</key>
        <false/>
        <!-- Restart if network becomes available -->
        <key>NetworkState</key>
        <true/>
    </dict>

    <!-- Throttle restarts (don't restart more than once per 10 seconds) -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Standard output/error logging -->
    <key>StandardOutPath</key>
    <string>/var/log/minnowvpn.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/minnowvpn.error.log</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/var/lib/minnowvpn</string>

    <!-- Run as root (required for TUN device creation) -->
    <key>UserName</key>
    <string>root</string>
    <key>GroupName</key>
    <string>wheel</string>

    <!-- File creation mask (stricter permissions on created files) -->
    <key>Umask</key>
    <integer>63</integer> <!-- 0077 in octal - owner only -->

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>RUST_LOG</key>
        <string>info</string>
        <key>RUST_BACKTRACE</key>
        <string>1</string>
    </dict>

    <!-- Process type for resource scheduling -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <!-- Nice priority (lower = higher priority, -5 for VPN responsiveness) -->
    <key>Nice</key>
    <integer>-5</integer>

    <!-- Timeout for graceful shutdown (seconds) -->
    <key>ExitTimeOut</key>
    <integer>30</integer>

    <!-- Enable sandbox (requires sandbox profile) -->
    <!-- Note: Uncomment when sandbox profile is tested and verified -->
    <!-- <key>SandboxProfile</key> -->
    <!-- <string>/Library/Application Support/MinnowVPN/com.minnowvpn.vpn-service.sb</string> -->

    <!-- Soft resource limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <!-- Max open files -->
        <key>NumberOfFiles</key>
        <integer>1024</integer>
    </dict>

    <!-- Hard resource limits -->
    <key>HardResourceLimits</key>
    <dict>
        <!-- Max open files -->
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>

    <!-- Session type - run in system session -->
    <key>LimitLoadToSessionType</key>
    <array>
        <string>System</string>
    </array>

    <!-- Enable mach service for on-demand activation -->
    <key>MachServices</key>
    <dict>
        <key>com.minnowvpn.vpn-service</key>
        <true/>
    </dict>

    <!-- Don't wait for debugger -->
    <key>WaitForDebugger</key>
    <false/>

    <!-- Abandon process group on termination -->
    <key>AbandonProcessGroup</key>
    <false/>
</dict>
</plist>
